# apps/api/src/sandbox-providers/google-run/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err

[program:socat]
command=socat TCP-LISTEN:2222,fork,reuseaddr TCP:localhost:22
autostart=true
autorestart=true
stdout_logfile=/var/log/socat.log
stderr_logfile=/var/log/socat.err

[program:nextjs]
command=/app/startup.sh
environment=PORT=3001
autostart=true
autorestart=true
stdout_logfile=/var/log/nextjs.log
stderr_logfile=/var/log/nextjs.err
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err
priority=20
depends_on=nextjs
startretries=10 